# raspberry_client/picms_player.service
# PiCMS Player Systemd Service File
# Location: /etc/systemd/system/picms_player.service

[Unit]
Description=PiCMS Video Player
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=video
SupplementaryGroups=render audio
WorkingDirectory=/home/pi/picms
ExecStart=/usr/bin/python3 /home/pi/picms/player.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# GPU/DRM access for video playback
Environment="XDG_RUNTIME_DIR=/run/user/1000"

[Install]
WantedBy=multi-user.target
